# gnu_linux/home/lg/crontab
# Clean up remaining ControlMaster sockets on reboot
@reboot rm -f ${HOME}/.ssh/ctl/*
