# gnu_linux/etc/init/lxdm.conf
# lxdm - LXDE Display Manager
#
# The display manager service manages the X servers running on the
# system, providing login and auto-login services
# Based on gdm upstart script
#
# Added "galaxy-ok" to start requirements for Liquid Galaxy

description	"GNOME Display Manager"
author		"William Jon McCann <mccann@jhu.edu>"

start on (filesystem
          and started dbus
          and galaxy-ok
          and (graphics-device-added fb0 PRIMARY_DEVICE_FOR_DISPLAY=1
               or drm-device-added card0 PRIMARY_DEVICE_FOR_DISPLAY=1
               or stopped udevtrigger))
stop on runlevel [016]

emits starting-dm

env XORGCONFIG=/etc/X11/xorg.conf

script
    if [ -n "$UPSTART_EVENTS" ]
    then
	[ ! -f /etc/X11/default-display-manager -o "$(cat /etc/X11/default-display-manager 2>/dev/null)" = "/usr/sbin/lxdm" ] || { stop; exit 0; }

	# Check kernel command-line for inhibitors
	for ARG in $(cat /proc/cmdline)
	do
	    case "${ARG}" in
		text|-s|s|S|single)
		    exit 0
		    ;;
	    esac
	done
    fi

    if [ -r /etc/default/locale ]; then
	. /etc/default/locale
	export LANG LANGUAGE
    elif [ -r /etc/environment ]; then
	. /etc/environment
	export LANG LANGUAGE
    fi

    #Export PATH from /etc/environment
    if [ -r /etc/environment ]; then
	. /etc/environment
	export PATH
    fi

    export XORGCONFIG

    exec lxdm-binary $* $CONFIG_FILE
end script
