# gnu_linux/etc/profile.d/lg-export.sh
# place in /etc/profile.d
[ -r ${HOME}/etc/shell.conf ] && . ${HOME}/etc/shell.conf && export FRAME_NO=$FRAME_NO
